# deploy/systemd/kalshi-gh1-websocket-drain.service
[Unit]
Description=Kalshi GH-1 guarded single-writer staging drain
After=network-online.target
ConditionPathExists=/etc/kalshi-bot/kalshi-bot.env

[Service]
Type=oneshot
User=kalshi
WorkingDirectory=/opt/kalshi-predictive-bot
EnvironmentFile=/etc/kalshi-bot/kalshi-bot.env
Environment=UI_READ_ONLY=true
Environment=EXECUTION_ENABLED=false
Environment=EXECUTION_DRY_RUN=true
Environment=EXECUTION_KILL_SWITCH=true
Environment=AUTOPILOT_DRY_RUN=true
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/kalshi-predictive-bot/.venv/bin/kalshi-bot gh1-websocket-orderbook-drain --apply
TimeoutStartSec=60s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryMax=384M
TasksMax=48
ReadWritePaths=/opt/kalshi-predictive-bot/reports/phase_gh1 /var/lib/kalshi-bot
ReadOnlyPaths=/opt/kalshi-predictive-bot /etc/kalshi-bot
