# deploy/systemd/kalshi-ui-status-collector.service
[Unit]
Description=Kalshi UI read-only status collector
After=network-online.target
ConditionPathExists=/etc/kalshi-bot/kalshi-bot.env

[Service]
Type=oneshot
User=kalshi
WorkingDirectory=/opt/kalshi-predictive-bot
EnvironmentFile=/etc/kalshi-bot/kalshi-bot.env
Environment=UI_READ_ONLY=true
Environment=EXECUTION_ENABLED=false
Environment=EXECUTION_DRY_RUN=true
Environment=EXECUTION_KILL_SWITCH=true
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/kalshi-predictive-bot/.venv/bin/python /opt/kalshi-predictive-bot/scripts/ui_obs3b_live_status_collector.py --output /opt/kalshi-predictive-bot/reports/ui_obs_live/progress_snapshot.json --backup-root /opt/kalshi-predictive-bot/reports/backup_metadata --service kalshi-r5-watcher.service
TimeoutStartSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryMax=192M
TasksMax=32
ReadWritePaths=/opt/kalshi-predictive-bot/reports/ui_obs_live
ReadOnlyPaths=/opt/kalshi-predictive-bot /etc/kalshi-bot

[Install]
WantedBy=multi-user.target
