# deploy/systemd/kalshi-gh2-decision-refresh.service
[Unit]
Description=Kalshi GH-2 guarded active-candidate decision refresh
After=network-online.target kalshi-gh1-websocket-watch.service
Wants=network-online.target
ConditionPathExists=/etc/kalshi-bot/kalshi-bot.env
ConditionPathExists=/var/lib/kalshi-bot/kalshi_phase1.db

[Service]
Type=oneshot
User=kalshi
WorkingDirectory=/opt/kalshi-predictive-bot
EnvironmentFile=/etc/kalshi-bot/kalshi-bot.env
Environment=UI_READ_ONLY=true
Environment=EXECUTION_ENABLED=false
Environment=EXECUTION_DRY_RUN=true
Environment=EXECUTION_KILL_SWITCH=true
Environment=AUTOPILOT_ENABLED=false
Environment=AUTOPILOT_DRY_RUN=true
Environment=KALSHI_WEBSOCKET_STAGING_DIR=/var/lib/kalshi-bot-gh1/staging
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/kalshi-predictive-bot/scripts/cloud/kalshi-gh2-decision-refresh.sh
TimeoutStartSec=6min
StateDirectory=kalshi-bot-gh2
StateDirectoryMode=0750
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryMax=768M
TasksMax=64
ReadWritePaths=/var/lib/kalshi-bot /var/lib/kalshi-bot-gh1 /var/lib/kalshi-bot-gh2 /opt/kalshi-predictive-bot/reports
ReadOnlyPaths=/opt/kalshi-predictive-bot /etc/kalshi-bot
